# src/modules/rtspios/filesystem/root/etc/supervisor/conf.d/rtsp-server.conf
[program:rtsp-server]
user=pi
command=/root/rtsp-server.sh
autostart=true
autorestart=true
